# configs/systemd/openclaw-agent.service.example
[Unit]
Description=Armando AI Assistant - OpenClaw Agent
After=network.target
Wants=network.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/home/openclaw/openclaw
ExecStart=/usr/bin/node /home/openclaw/openclaw/dist/index.js gateway
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=OPENCLAW_CONFIG_DIR=/home/openclaw/.openclaw
Environment=TELEGRAM_BOT_TOKEN=${TELEGRAM_BOT_TOKEN}
Environment=DEEPSEEK_API_KEY=${DEEPSEEK_API_KEY}

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/openclaw/.openclaw /var/log/openclaw

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-agent

[Install]
WantedBy=multi-user.target
