# configs/systemd/openclaw-agent.service.example
[Unit]
Description=Armando AI Assistant
After=network.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/openclaw
ExecStart=/usr/bin/node /home/openclaw/openclaw/dist/index.js gateway
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
